# getSun.cron
0 0 * * * root /usr/bin/python /home/pi/Circadian-Lights/getSun.py
