# sunRise.cron
50 07 * * * /usr/bin/python /home/pi/Circadian-Lights/sunrise.py