# sunrise.cron
50 07 * * * pi /usr/bin/python /home/pi/Circadian-Lights/sunrise.py