# sunset.cron
52 16 * * * pi /usr/bin/python /home/pi/Circadian-Lights/sunset.py